# deploy/systemd/vaam.service
[Unit]
Description=Vaam Gunicorn Daemon
After=network.target postgresql.service

[Service]
User=vaam
Group=www-data
WorkingDirectory=/home/vaam/app
EnvironmentFile=/home/vaam/app/.env
ExecStart=/home/vaam/app/venv/bin/gunicorn --config gunicorn.conf.py vaam_project.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
